# com.agent.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agent.monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
        <string>-u</string>
        <string>/Users/yvonne/Documents/Agent/main.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/yvonne/Documents/Agent/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yvonne/Documents/Agent/logs/launchd_stderr.log</string>
</dict>
</plist>
